# agent/rootfs/etc/s6-overlay/s6-rc.d/svc-desktop/run
#!/command/with-contenv bash

export DISPLAY=:0
export HOME=/home/claworc

# Wait for X server to be ready
while ! xdpyinfo -display :0 >/dev/null 2>&1; do sleep 0.5; done

# Start openbox window manager in background
s6-setuidgid claworc openbox &

# Run chromium in a restart loop
while true; do
  s6-setuidgid claworc /usr/bin/chromium \
    --no-first-run \
    --no-sandbox \
    --password-store=basic \
    --simulate-outdated-no-au='Tue, 31 Dec 2099 23:59:59 GMT' \
    --start-maximized \
    --user-data-dir=/home/claworc/chrome-data \
    --remote-debugging-port=9222 \
    --remote-debugging-address=127.0.0.1 \
    --disable-default-apps \
    --disable-features=CloseWindowWithLastTab > /dev/null 2>&1
  sleep 1
done
